# MaxFansService.service
[Unit]
Description=example systemd service unit file.

[Service]
ExecStart=/bin/bash ~/vastinstall/MaxFans.sh

[Install]
WantedBy=multi-user.target
